# contrib/launchd/dev.klb.suffuse.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.klb.suffuse.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/suffuse</string>
        <string>server</string>
    </array>

    <!-- Config file: ~/Library/Application Support/suffuse/suffuse.toml
         or set env vars below. -->
    <key>EnvironmentVariables</key>
    <dict>
        <!-- <key>SUFFUSE_TOKEN</key><string>mysecret</string> -->
        <!-- <key>SUFFUSE_ADDR</key><string>:9876</string> -->
    </dict>

    <!-- Logs go to ~/Library/Logs/suffuse/ -->
    <key>StandardOutPath</key>
    <string>/Users/Shared/suffuse-server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Shared/suffuse-server.log</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Throttle restart on repeated failures -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
